# private_dot_config/systemd/user/newsboat.timer
[Unit]
Description=Newsboat RSS Feed Update Daemon Timer

[Timer]
OnActiveSec=300
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
