# private_dot_config/systemd/user/newsboat.service
[Unit]
Description=Newsboat RSS Feed Update Daemon
Wants=newsboat.timer

[Service]
ExecStart=/snap/bin/newsboat -x reload

[Install]
WantedBy=default.target
